# orq-upstart.conf
description "orq daemon"
start on started docker
stop on stopping docker
respawn
exec /root/orq daemon
